# service/kaidadb.service
[Unit]
Description=KaidaDB Media Database Server
After=network.target

[Service]
Type=simple
ExecStart=@@INSTALL_DIR@@/kaidadb-server --config @@CONFIG_DIR@@/config.toml
Restart=always
RestartSec=5
LimitNOFILE=65536
Environment=RUST_LOG=kaidadb=info
WorkingDirectory=@@DATA_DIR@@

[Install]
WantedBy=default.target
